# scripts/systemd/sumbur-watch.path
[Path]
PathModified=/srv/web/sumbur.info/sumbur/sumbur

[Install]
WantedBy=multi-user.target
